# supervisord.conf
[supervisord]
nodaemon=true

[program:laravel]
command=php artisan serve --host=0.0.0.0 --port=80
autostart=true
autorestart=true
user=sail
redirect_stderr=true
stdout_logfile=/var/log/laravel.log
